# deploy/systemd/binance-hunter-monitor.service
[Unit]
Description=Binance pump-dump hunter monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/binance-hunter/backend
EnvironmentFile=/etc/binance-hunter.env
ExecStart=/opt/binance-hunter/backend/.venv/bin/python run.py monitor --config config/settings.json --top 120 --broad-top 450 --discover-every 15m --max-workers 3
Restart=always
RestartSec=5
User=root
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
