# deploy/systemd/binance-hunter-micro.service
[Unit]
Description=Binance hunter micro data collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/binance-hunter/backend
EnvironmentFile=/etc/binance-hunter.env
ExecStart=/opt/binance-hunter/backend/.venv/bin/python run.py collect-micro --config config/settings.json --broad-top 450 --depth-top 60 --depth-interval 30
Restart=always
RestartSec=10
User=root
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
